# tennis-bot.service
[Unit]
Description=Tennis Court Booking Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tennis_booking_bot
ExecStart=/usr/bin/python3 /home/ubuntu/tennis_booking_bot/src/telegram_bot.py
Restart=always
RestartSec=10

# Environment
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/home/ubuntu/tennis_booking_bot/.env

# Logging
StandardOutput=append:/home/ubuntu/tennis_booking_bot/logs/systemd.log
StandardError=append:/home/ubuntu/tennis_booking_bot/logs/systemd-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
